# cyberdino-panel.service
[Unit]
Description=CyberDino Painel Secundario (Xvfb + VNC)
After=network.target

[Service]
Type=simple
User=evo
ExecStart=/home/evo/linuxcnc/configs/Dino_Evo/start_panel.sh
ExecStop=/home/evo/linuxcnc/configs/Dino_Evo/start_panel.sh stop
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
